# deploy/tewodros-terminal.service
[Unit]
Description=tewodros-terminal portfolio
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/tewodros-terminal
ExecStart=/opt/tewodros-terminal/tewodros-terminal
Restart=always
RestartSec=5

Environment=SSH_HOST=0.0.0.0
Environment=SSH_PORT=22
Environment=HTTP_HOST=0.0.0.0
Environment=HTTP_PORT=443
Environment=TLS_CERT=/opt/tewodros-terminal/origin.crt
Environment=TLS_KEY=/opt/tewodros-terminal/origin.key
Environment=HOST_KEY_DIR=/opt/tewodros-terminal/.ssh
Environment=DB_PATH=/opt/tewodros-terminal/guestbook.db

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tewodros-terminal
PrivateTmp=true

AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
